# provisioning/roles/gunicorn/templates/gunicorn.conf.j2
description "gunicorn server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec {{ item.vars.command }}
